# deploy/padel-bot.service
[Unit]
Description=Padel Bot Web UI
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=padelbot
Group=padelbot
WorkingDirectory=/opt/padel-bot
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/padel-bot/.venv/bin/python -m uvicorn app.main:app --host 127.0.0.1 --port 18018
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
