# deploy/padel-bot.service
[Unit]
Description=Padel Bot Web UI
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=padelbot
Group=padelbot
WorkingDirectory=/opt/padel-bot
Environment=PYTHONUNBUFFERED=1
Environment=HOME=/home/padelbot
Environment=XDG_RUNTIME_DIR=/tmp/padelbot-runtime
Environment=DBUS_SESSION_BUS_ADDRESS=
ExecStartPre=/bin/mkdir -p /tmp/padelbot-runtime
ExecStartPre=/bin/chown padelbot:padelbot /tmp/padelbot-runtime
ExecStartPre=/bin/chmod 700 /tmp/padelbot-runtime
ExecStart=/usr/bin/xvfb-run -a --server-args=-screen\ 0\ 1280x900x24 /opt/padel-bot/.venv/bin/python -m uvicorn app.main:app --host 127.0.0.1 --port 18018
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
